# tutorial/package/supervisor/ztp-tftp.conf
[program:ztp-tftp]
command=/opt/ztp/venv/bin/python ztp_tftp.py
numprocs=1
directory=/opt/ztp/
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true