# tutorial/package/supervisor/ztp-tasks.conf
[program:ztp-tasks]
command=/opt/ztp/venv/bin/rq worker ztp_tasks
numprocs=1
directory=/opt/ztp/
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true